# debian-bookworm/debian/systemd/dhcpcd.service.d/quiet.conf
# Start dhcpcd in quiet mode to avoid double logging
[Service]
ExecStart=
ExecStart=/usr/sbin/dhcpcd -q -b
