# debian-bookworm/systemd/NetworkManager.service.d/restart.conf
# Set NetworkManager to auto-restart since it's critical infrastructure

[Service]
Restart=always
RestartSec=30s
